# dot_config/cron.daily/executable_xlocate
#!/bin/sh
export PATH="$HOME/void/xtools:$PATH"
command -v xlocate >/dev/null 2>&1 && exec xlocate -S
